# src/usr/lib/systemd/system/wslfree.service
[Unit]
Description = Free WSL2 kernel caches

[Service]
Type = oneshot
ExecStart = /usr/bin/wslfree

[Install]
WantedBy = multi-user.target
