# src/usr/lib/systemd/system/wslfree.timer
[Unit]
Description = Timer to free WSL2 kernel caches

[Timer]
OnBootSec=1h
OnUnitActiveSec=1h

[Install]
WantedBy = timers.target
